# cron.d/jenkins-docker
# Note: It's safest to copy this file into /etc/cron.d rather than
# symlinking it as cron is strict about ownership of the cron.d files.
@reboot jenkins [ -x /store/scripts/cron.d/jenkins-docker.sh ] && /store/scripts/cron.d/jenkins-docker.sh
